# services/project-runner/project-runner.service
[Unit]
Description=OrquestrAI project-runner (clone daemon, Fase A do import)
# docker.service: o bind e no IP da bridge app-net (172.18.0.1) — precisa
# do docker de pe. Restart=on-failure cobre a corrida se ainda assim falhar.
After=network.target docker.service

[Service]
Type=simple
User=projrunner
Group=projrunner
EnvironmentFile=/etc/project-runner.env
ExecStart=/usr/bin/node /opt/project-runner/server.js
Restart=on-failure
RestartSec=3

# --- Hardening systemd (defesa em profundidade, mesmo sendo nao-root) ---
NoNewPrivileges=yes
ProtectSystem=strict
# unico caminho gravavel: o staging. Nada mais do host e escrivel.
ReadWritePaths=/var/www/orquestrai/projects/.staging
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictSUIDSGID=yes
LockPersonality=yes
# git precisa de rede (https) e o server escuta TCP: so IPv4/IPv6
RestrictAddressFamilies=AF_INET AF_INET6
# sem capabilities (nao-root ja, mas explicito)
CapabilityBoundingSet=
AmbientCapabilities=
# limites de recurso: clone raso e barato
MemoryMax=256M
TasksMax=64
# arquivos do clone nascem 640/750 (grupo projrunner, sem world-read)
UMask=0027

[Install]
WantedBy=multi-user.target
